# examples/systemd/backup.service
[Unit]
Description=Backup Starter

[Service]
ExecStart=/usr/bin/backup -loglevel=debug -config=/etc/backup/config.yaml